# files/i3/systemd_setI3Config.service
[Unit]
Description=Set i3 symlinks

[Service]
Type=oneshot
ExecStart=~/.config/i3/set_i3_config_simlinks.sh
[Install]
WantedBy=multi-user.target